# selfhosted/resources/restic.timer
[Unit]
Description=Make a daily backup for %I

[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=timers.target
